# mcannounce-nft.service
[Service]
Type=simple
ExecStart=/usr/bin/mcannounce-nft inet filter mcport
SyslogIdentifier=%N

[Install]
WantedBy=multi-user.target
